# .config/systemd/user/hydrate.timer
[Unit]
Description=Hourly reminder to drink water

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
