# .config/systemd/user/pass-push.service
[Unit]
Description=Push pass git store to remote

[Service]
Type=oneshot
WorkingDirectory=%h/.password-store
ExecStart=/usr/bin/git push
