# .config/systemd/user/pass-push.timer
[Unit]
Description=Daily git push for pass

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
